# home/.config/autostart/caps_remap.desktop
[Desktop Entry]
Categories=System;
Exec=caps_remap.sh
Name=gebaard
StartupNotify=false
Terminal=false
Type=Application
