# home/hidpi/kde/.config/autostart/gebaard.desktop
[Desktop Entry]
Categories=System;
Exec=gebaard -b
Name=gebaard
StartupNotify=false
Terminal=false
Type=Application
